# openwrt/netforensics/files/forensicsd.init
#!/bin/sh /etc/rc.common
# procd init for forensicsd

START=95
USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/sbin/forensicsd --netlink \
		--router-id "$(cat /proc/sys/kernel/hostname 2>/dev/null || echo cpe)" \
		--wifi-if wlan0 --wifi-interval-ms 10000
	procd_set_param respawn 3600 5 0
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
}
