# deploy/cpe_agent.service
[Unit]
Description=netforensics CPE agent daemon (perf/tcp NDJSON → edgehost)
Documentation=file:///opt/netforensics/docs/guides/cpe-agent-daemon.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cpe_agent --config /etc/cpe_agent/cpe_agent.yaml --router-id %H
ExecReload=/bin/kill -HUP $MAINPID
# Live ICMP: CAP_NET_RAW; NFLOG TCP stats: CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
MemoryMax=32M
TasksMax=16
LimitNOFILE=256
ReadWritePaths=/var/spool/netforensics /var/run/netforensics
RuntimeDirectory=netforensics
RuntimeDirectoryMode=0755
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
