# deploy/cpe_agent.service
[Unit]
Description=netforensics CPE performance agent (cpe_perf NDJSON)
Documentation=file:///opt/netforensics/ARCHITECTURE.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cpe_agent --config /etc/cpe_agent/cpe_agent.yaml --router-id %H
# Live ICMP may need CAP_NET_RAW when SOCK_DGRAM IPPROTO_ICMP is unavailable.
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
MemoryMax=24M
TasksMax=16
LimitNOFILE=128
ReadWritePaths=/var/spool/netforensics
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
